# dotfiles/.config/systemd/user/system-health-review.timer
[Unit]
Description=Daily system health review

[Timer]
OnCalendar=*-*-* 09:00:00
AccuracySec=1h
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
